# debian/pdns-backend-sqlite.prerm
#!/bin/sh
#
# Pre removal

set -e

# Stop pdns.
if [ -x "/etc/init.d/pdns" ]; then
	invoke-rc.d pdns stop || exit $?
fi

# Execute dbconfig-common
. /usr/share/debconf/confmodule
. /usr/share/dbconfig-common/dpkg/prerm.sqlite
dbc_go pdns-backend-sqlite $@

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0

