# debian/pdns-backend-lua.prerm
#!/bin/sh
#
# Add call to init script to actually stop pdns.

if [ -x "/etc/init.d/pdns" ]; then
	invoke-rc.d pdns stop || exit $?
fi

exit 0

